# ota_sync/deploy/amap-ota-sync.timer
[Unit]
Description=Periodically sync AMap ESP32 OTA artifacts

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Unit=amap-ota-sync.service
Persistent=true

[Install]
WantedBy=timers.target
